# Python/scheduler/crontab
# Transfer Enigma Scheduler
# Format: minute hour day month weekday command
# Output is redirected to /proc/1/fd/1 (Docker stdout) so it appears in container logs.

PYTHONPATH=/apps
HOME=/apps

# Fetch exchange rates from CBR every day at 06:00
* * * * * python /apps/scheduler/run_job.py --job-name update_rates > /proc/1/fd/1 2>/proc/1/fd/2
